# build/packer/scripts/rxe.service
[Unit]
Description=SoftRoCE
Requires=rdma-hw.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/rxe_all.sh
ExecStop=/usr/bin/rxe_cfg stop

[Install]
WantedBy=multi-user.target
